# daemons/metamanager-compress-daemon.service
[Unit]
Description=Metamanager Lossless Compression Daemon
Documentation=https://github.com/richardkentgates/metamanager
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/metamanager-compress-daemon.sh
User=www-data
Group=www-data
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ReadWritePaths=__WP_CONTENT_DIR__/metamanager-jobs __WP_CONTENT_DIR__/uploads /tmp /var/log

# Logging to journald (also appended to log file by the script itself)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=metamanager-compress

[Install]
WantedBy=multi-user.target
